# etc/profile.d/pythonpath.sh
export PYTHONPATH=${PYTHONPATH}:/usr/local/lib/python3/dist-packages
